# data_backup/backup.cron
# 
1 1 * * * /work/yonder/data_backup/sh/backup.sh

# end line